# .config/systemd/user/nvidia-display-monitor.service
[Unit]
Description=NVIDIA Display Monitor - Auto-apply layout on display changes
After=graphical-session.target

[Service]
Type=simple
ExecStart=%h/.screenlayout/display-monitor.sh
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
